# templates/etc/systemd/system/hhvm.service.erb
[Unit]
Description=HipHop Virtual Machine

[Service]
LimitNOFILE=><%= @limit %>
LimitNPROC=><%= @limit %>
RuntimeDirectory=hhvm
ExecStartPre=/bin/chown <%= scope.lookupvar('hhvm::params::hhvm_user') %>:<%= scope.lookupvar('hhvm::params::hhvm_group') %> /var/run/hhvm
ExecStart=<%= scope.lookupvar('hhvm::path_to_hhvm') %> -m daemon -u <%= scope.lookupvar('hhvm::params::hhvm_user') %> -c <%= @server_ini %> -c <%= @php_ini %> -c <%= @config_hdf %>
ExecStop=rm <%= @pid %>
PIDFile=<%= @pid %>
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
