# templates/etc/init/hhvm.conf.erb
# hhvm - HipHop VM
#
# The HipHopVM server provides a high performance PHP stack and web server.

description     "HHVM server"
author "Yermo Lamers http://twitter.com/yermolamers"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022
limit nofile <%= scope.lookupvar('hhvm::limit') %> <%= scope.lookupvar('hhvm::limit') %>

post-stop script
    /bin/rm -f <%= scope.lookupvar('hhvm::params::pid') %> <%= scope.lookupvar('hhvm::params::socket') %>
end script

# Location of config and executable

env SERVER=<%= scope.lookupvar('hhvm::path_to_hhvm') %>
env CONFIG_FILE="/etc/hhvm/server.ini"
env SYSTEM_CONFIG_FILE="/etc/hhvm/php.ini"
env HDF_CONFIG_FILE="/etc/hhvm/config.hdf"
env RUN_AS=www-data
exec $SERVER -m daemon --config $CONFIG_FILE --config $SYSTEM_CONFIG_FILE --config $HDF_CONFIG_FILE --user $RUN_AS
